# voice-prompt.service
[Unit]
Description=Voice Prompt — voice-to-clipboard prompt generator
Documentation=https://github.com/tr4m0ryp/voice-prompt
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/voice-prompt
Restart=on-failure
RestartSec=3
Environment=DISPLAY=:0
Environment=XDG_SESSION_TYPE=wayland

[Install]
WantedBy=default.target
